# services/monolith/ansible/roles/mirror-hdd/templates/mirror-hdd.timer.j2
[Unit]
Description=Nightly mirror of hdd-c to hdd-d
Requires=mirror-hdd.service

[Timer]
OnCalendar={{ mirror_on_calendar }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
