# services/watchtower/ansible/roles/reolink_exporter/templates/reolink_exporter.service.j2
[Unit]
Description=Reolink NVR Prometheus Exporter
After=network.target

[Service]
Type=simple
User=reolink_exporter
Group=reolink_exporter
ExecStart=/usr/bin/python3 /opt/reolink_exporter/reolink_exporter.py
Restart=on-failure
RestartSec=10

Environment=NVR_HOST={{ reolink_nvr_host }}
Environment=NVR_USER={{ reolink_nvr_user }}
Environment=NVR_PASSWORD={{ vault_reolink_nvr_password }}
Environment=SCRAPE_PORT={{ reolink_scrape_port }}
Environment=POLL_INTERVAL={{ reolink_poll_interval }}

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
